# supervisord.conf
[program:node]
command = node game.js
directory = /home/dotcloud/current
